# supervisord.conf
[supervisord]
nodaemon=true

[program:worker_user_dxf_2_svg]
command=python worker_user_dxf_2_svg.py
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr

